# ops/systemd/bureau-status-projection.service
[Unit]
Description=Bureau read-only status projection board
Documentation=https://github.com/heimgewebe/bureau/blob/main/docs/bureau-runtime-observation-v1.md
ConditionPathIsDirectory=%h/repos/bureau
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/bureau/venv/bin/bureau --root %h/repos/bureau --json status-projection
Environment=PYTHONUNBUFFERED=1
Environment=GH_NO_UPDATE_NOTIFIER=1
Environment=GH_PROMPT_DISABLED=1
TimeoutStartSec=3min
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077
